# dist/DEBIAN/postinst
#!/bin/sh
systemctl enable fxtrigger.service
systemctl start fxtrigger.service
