# dist/DEBIAN/prerm
#!/bin/sh
systemctl stop fxtrigger.service
systemctl disable fxtrigger.service
